# jobs/nro-extractor/crontab
# m   h  dom mon dow user    command
  */2 *   *   *   *  default /code/run.sh 
